# gps-pi/rc.local
#!/bin/bash
echo -e "$(cat /proc/uptime | awk '{print $1}') Seconds" > /var/log/boottime
if (( $(cat /DietPi/dietpi/.install_stage) == 1 )); then

    /DietPi/dietpi/dietpi-services start

fi
/DietPi/dietpi/dietpi-banner 0
/root/gps-pi/log.sh &
echo -e " Default Login:\n Username = root\n Password = dietpi\n"
exit 0
